# auto_down.cron
*/1 * * * * /usr/src/django_test/auto_down.sh